# com.tradingcontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingcontrol.v3system</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/matthew/Desktop/trading-control-python && ./boot_v3.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/matthew/Desktop/trading-control-python</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/matthew/Desktop/trading-control-python/logs/v3_system.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matthew/Desktop/trading-control-python/logs/v3_system_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/matthew/Desktop/trading-control-python</string>
    </dict>
</dict>
</plist>
